# root/etc/uci-defaults/40_luci-tailscale
#!/bin/sh

UPGRADE_STATE_DIR=/etc/.luci-app-tailscale-upgrade

cleanup_upgrade_state() {
	rm -f "$UPGRADE_STATE_DIR/tailscale" "$UPGRADE_STATE_DIR/tailscale_openclash" \
		"$UPGRADE_STATE_DIR/tailscale_openclash.absent" "$UPGRADE_STATE_DIR/.complete"
	rmdir "$UPGRADE_STATE_DIR" 2>/dev/null || true
}

restore_upgrade_config() {
	local source_file="$1" target_file="$2" temporary_file

	[ -f "$source_file" ] || return 0
	temporary_file="$target_file.luci-app-tailscale-upgrade.$$"
	cp "$source_file" "$temporary_file" || return 1
	mv "$temporary_file" "$target_file" || {
		rm -f "$temporary_file"
		return 1
	}
}

if [ -f "$UPGRADE_STATE_DIR/.complete" ]; then
	restore_upgrade_config "$UPGRADE_STATE_DIR/tailscale" /etc/config/tailscale || exit 1
	if [ -f "$UPGRADE_STATE_DIR/tailscale_openclash.absent" ]; then
		rm -f /etc/config/tailscale_openclash
	else
		restore_upgrade_config "$UPGRADE_STATE_DIR/tailscale_openclash" /etc/config/tailscale_openclash || exit 1
	fi
	cleanup_upgrade_state
elif [ -d "$UPGRADE_STATE_DIR" ]; then
	cleanup_upgrade_state
fi

/usr/sbin/tailscale_secrets migrate || exit 1

openclash_bypass_enabled="$(uci -q get tailscale_openclash.settings.enabled 2>/dev/null || true)"
if [ "$openclash_bypass_enabled" = '1' ] && [ -x /etc/init.d/tailscale-openclash-bypass ]; then
	/etc/init.d/tailscale-openclash-bypass enable
	/etc/init.d/tailscale-openclash-bypass start
fi

rm -f /tmp/luci-indexcache
exit 0
